# coreOS-configs/netbox/quadlets/netbox-redis.container
[Unit]
Description=NetBox Redis
# Stay inert until the operator creates the real env file (copy redis.env.example
# → redis.env and fill in the password). Avoids starting with a placeholder.
ConditionPathExists=%h/netbox/env/redis.env

[Container]
ContainerName=netbox-redis
Image=docker.io/valkey/valkey:9.0-alpine
Pod=netbox.pod
EnvironmentFile=%h/netbox/env/redis.env
Volume=netbox-redis-data:/data
Exec=sh -c 'valkey-server --appendonly yes --requirepass "$REDIS_PASSWORD"'

[Service]
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=default.target
