# coreOS-configs/netbox/quadlets/netbox-postgres.container
[Unit]
Description=NetBox Postgres
# Stay inert until the operator creates the real env file (copy postgres.env.example
# → postgres.env and fill in the password). Avoids initializing with a placeholder.
ConditionPathExists=%h/netbox/env/postgres.env

[Container]
ContainerName=netbox-postgres
Image=docker.io/postgres:18-alpine
Pod=netbox.pod
EnvironmentFile=%h/netbox/env/postgres.env
Environment=POSTGRES_DB=netbox
Environment=POSTGRES_USER=netbox
Volume=netbox-postgres:/var/lib/postgresql

[Service]
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=default.target
