# coreOS-configs/netbox/quadlets/netbox-worker.container
[Unit]
Description=NetBox Worker
After=netbox.service
Requires=netbox.service
# Shares netbox.env with the netbox service; stay inert until it exists.
ConditionPathExists=%h/netbox/env/netbox.env

[Container]
ContainerName=netbox-worker
Image=docker.io/netboxcommunity/netbox:v4.5-4.0.2
Pod=netbox.pod
User=netbox
Group=root
EnvironmentFile=%h/netbox/env/netbox.env
Environment=DB_USER=netbox
Environment=REDIS_CACHE_PORT=6380
Environment=CORS_ORIGIN_ALLOW_ALL=True
Volume=%h/netbox/config:/etc/netbox/config:z,ro
Volume=netbox-media-files:/opt/netbox/netbox/media:rw
Volume=netbox-reports-files:/opt/netbox/netbox/reports:rw
Volume=netbox-scripts-files:/opt/netbox/netbox/scripts:rw
Exec=/opt/netbox/venv/bin/python /opt/netbox/netbox/manage.py rqworker
HealthCmd=sh -c "ps -aux | grep -v grep | grep -q rqworker || exit 1"
HealthStartPeriod=20s
HealthInterval=15s
HealthTimeout=3s

[Service]
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=default.target
